# 0x1A-application_server/gunicorn.service
[Unit]
Description=gunicorn daemon for hbnb
After=network.target
Requires=gunicorn.socket

[Service]
User=ubunt
Group=www-data1

WorkingDirectory=/home/ubuntu/AirBnB_clone_v4
ExecStart=/usr/bin/gunicorn \
          --error-logfile /tmp/airbnb-error.log \
          --access-logfile /tmp/airbnb-access.log \
          --bind unix:/0.0.0.0:5003 \
          --workers 3 \
          web_dynamic.2-hbnb:app

[Install]
WantedBy=multi-user.target
